# systemd/kingdom4-commander.service
[Unit]
Description=Kingdom 4 Commander
After=network.target

[Service]
WorkingDirectory=/path/to/kingdom4
ExecStart=/path/to/kingdom4/venv/bin/uvicorn commander.commander:app --host 0.0.0.0 --port 8000
Restart=always

[Install]
WantedBy=multi-user.target
